# etc/systemd/system/kiosk.service
[Unit]
Description=TC8 fullscreen kiosk (weston + browser)
After=systemd-user-sessions.service seatd.service tc8-persist-root.service
Wants=seatd.service
ConditionPathExists=/dev/dri/card0

[Service]
Type=simple
User=kiosk
Group=kiosk
PAMName=login
TTYPath=/dev/tty7
TTYReset=no
TTYVHangup=no
StandardInput=tty
StandardOutput=journal
StandardError=journal
EnvironmentFile=-/etc/default/tc8-kiosk
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=XDG_SEAT=seat0
Environment=HOME=/home/kiosk
Environment=XDG_SESSION_TYPE=wayland
Environment=XDG_CACHE_HOME=/persist/kiosk/cache
Environment=XDG_DATA_HOME=/persist/kiosk/share
ExecStartPre=+/bin/sh -c "install -d -o kiosk -g kiosk -m 0700 /run/user/1000 /persist/kiosk/cache /persist/kiosk/share || true"
ExecStart=/usr/local/bin/kiosk-launch
Restart=always
RestartSec=3

[Install]
WantedBy=graphical.target
