# etc/systemd/system/tc8-config.service
[Unit]
Description=Apply TC8 config blob pushed to the cache partition (provisioning wizard)
Documentation=file:///etc/tc8-config/apply-config.sh
# Needs /dev/disk/by-partlabel + a writable /etc; runs before the kiosk reads its
# env. Ordered before kiosk-config.service so a local /data override still wins.
After=local-fs.target tc8-persist-root.service
Before=systemd-networkd.service wpa_supplicant@wlan0.service kiosk-config.service kiosk.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/etc/tc8-config/apply-config.sh

[Install]
WantedBy=multi-user.target
